# deployment/systemd/aeon-stack.service
[Unit]
Description=AEON stack (docker compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/srv/stack
ExecStart=/usr/bin/docker compose -f deployment/docker-compose.production.yml up -d
ExecStop=/usr/bin/docker compose -f deployment/docker-compose.production.yml down
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target 